# conf/supervisor.conf
# -*- conf -*-

[program:brevity]
command=/srv/sites/brevity/conf/startup
user=ebower
group=web
autostart=true
autorestart=true
stderr_logfile=syslog
stdout_logfile=syslog
